# rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
eventtrack.sh &
moduleinit.sh &
wandete.sh &
/home/sensegrow/monitor.sh &
/home/sensegrow/gpio.sh &
/bin/date +%s -s @`cat /overlay/home/sensegrow/last_reboot_timestamp` &
exit 0
